# provisioning/supervisor/mobovidata_dj.conf
[program:mobovidata_dj]
command = /www_app_group/mobovidata/bin/gunicorn_start                    ; Command to start app
user = www_user                                                          ; User to run as
stdout_logfile = /www_app_group/mobovidata/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
